# posthog-telegram.service
[Unit]
Description=PostHog → Telegram event bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/posthog-telegram
ExecStart=/opt/posthog-telegram/venv/bin/python bot.py
EnvironmentFile=/opt/posthog-telegram/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
